# cpcbf/field/cpcbf-lora.service
[Unit]
Description=CPCBF Autonomous Benchmark — LoRa Track (MKR WAN 1300)
After=multi-user.target bluetooth.service
Wants=multi-user.target bluetooth.service

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 __SCRIPT_DIR__/auto_bench.py __SCRIPT_DIR__/role_lora.json
StandardOutput=journal
StandardError=journal
RemainAfterExit=yes
TimeoutStartSec=infinity

[Install]
WantedBy=multi-user.target
